# examples/com.agent-chat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agent-chat</string>

  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/python3</string>
    <string>/ABSOLUTE/PATH/TO/agent-chat/agent_chat_control_plane.py</string>
    <string>run</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>AGENT_IMESSAGE_TO</key>
    <string>+15555550123</string>

    <key>AGENT_CHAT_HOME</key>
    <string>/Users/USERNAME/.codex</string>

    <key>AGENT_CHAT_NOTIFY_MODE</key>
    <string>route</string>

    <key>AGENT_CHAT_LAUNCHD_LABEL</key>
    <string>com.agent-chat</string>

    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/USERNAME/Library/Logs/agent-chat.launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/USERNAME/Library/Logs/agent-chat.launchd.err.log</string>
</dict>
</plist>
